# lib/cmbackup/cmbackup.cron
SHELL=/bin/bash
PATH=/sbin:/bin:/usr/sbin:/usr/bin:/usr/local/bin
MAILTO=root

# Every 1st Sunday of the month
01 0 * * 7     zextras [ $(date +\%d) -le 8 ] && ( /opt/cbplus/utils/cmbackup_full )
10 1 * * *     zextras /opt/cbplus/utils/cmbackup_inc

